# build/ringrem/DEBIAN/prerm
#!/bin/sh
set -e

systemctl stop ringrem.timer
systemctl disable ringrem.timer
systemctl daemon-reload

exit 0

