# build/ringrem/systemd/system/ringrem.service
[Unit]
Description=Ringrem notification check

[Service]
Type=oneshot
ExecStart=/usr/bin/ringrem notify
